# crontab
@reboot sudo python3 /Desktop/summer/Python/SiteBlocker/block.py
